# YourFlix Setup/YourFlixMKIII_A1.15.D1.14_all/lib/systemd/system/yourflix.service
[Unit]
Description=Autoscan YourFlix Scanner
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3 /usr/share/yourflix/yf_scanner.py
Restart=always

[Install]
WantedBy=multi-user.target